# private_dot_config/systemd/user/cleanup.service
[Unit]
Description=Daily cleanup tasks

[Service]
Type=oneshot
ExecStart=%h/.local/bin/empty-downloads 3
ExecStart=%h/.local/bin/empty-trash 7
