# private_dot_config/systemd/user/empty-trash.service
[Unit]
Description=Empty trash

[Service]
Type=oneshot
ExecStart=%h/.local/bin/empty-trash
